# dot_config/systemd/user/agent-url-opener.service
[Unit]
Description=Open agent-isolated sandbox URLs in the host browser

[Service]
Type=oneshot
ExecStart=%h/.local/libexec/agent-url-drain
