# dot_config/systemd/user/agent-url-opener.path
[Unit]
Description=Watch agent-isolated sandbox URL spool and open links in the host browser

[Path]
# Fires agent-url-opener.service whenever the sandbox drops a *.url file here.
PathExistsGlob=%h/.local/state/agent-url-open/*.url

[Install]
WantedBy=default.target
