# dot_config/systemd/user/dropdown-terminal.service
[Unit]
Description=Dropdown terminal (alacritty)
PartOf=graphical-session.target
After=graphical-session.target

[Service]
# Spawn alacritty with a fixed app_id so the windowrules in
# dot_config/hypr/windows.conf can pin it to special:dropdown. Unlike
# Ghostty (which silently dropped any --class that wasn't a valid GTK
# reverse-DNS app ID), alacritty sets the Wayland app_id to the first
# `--class` value verbatim — the `com.mitchellh.ghostty-dropdown` string
# is kept only because it's the established match token in windows.conf;
# renaming it would mean touching every matching rule for no functional gain.
#
# Restart=always keeps the workspace pre-warmed — closing the dropdown
# respawns within 1s so the next $mod+u summon is instant.
ExecStart=/usr/bin/alacritty --class=com.mitchellh.ghostty-dropdown
Restart=always
RestartSec=1
# Cap cascading restarts (alacritty crash-loop after a config edit) — fail
# loud after 5 retries in 60s instead of burning CPU silently.
StartLimitBurst=5
StartLimitIntervalSec=60

[Install]
WantedBy=graphical-session.target
