# dot_config/systemd/user/helium-update.timer
[Unit]
Description=Weekly Helium AppImage update check

[Timer]
OnBootSec=10min
OnUnitActiveSec=1w
Persistent=true
RandomizedDelaySec=1h

[Install]
WantedBy=timers.target
