# dot_config/systemd/user/hyprwhspr-rewrite-e4b.service
[Unit]
Description=hyprwhspr dictation rewrite backend — Gemma-4 E4B + MTP (llama.cpp Vulkan, :8083)
Documentation=https://github.com/NullSense/hyprwhspr-ai
After=graphical-session.target

[Service]
Type=simple
# Gemma-4 E4B Q5_K_M, single-stream (NP=1), summarization sampler (thinking off),
# MTP speculative decoding ON (E4B assistant draft head, ~80% accept → ~40% faster
# dictation rewrites, lossless). Bench: tests/rewrite-bench p50 136ms vs 227ms plain.
Environment=PORT=8083
Environment=NP=1
Environment=MODE=sum
Environment=SPEC=on
ExecStart=%h/bin/llm-servers/llama-gemma4b
Restart=on-failure
RestartSec=3
TimeoutStopSec=10

[Install]
WantedBy=default.target
