# dot_config/systemd/user/hyprwhspr-ai.service
[Unit]
Description=hyprwhspr-ai coordinator daemon (rewrite + vision + translate)
Documentation=https://github.com/goodroot/hyprwhspr
After=graphical-session.target

[Service]
Type=simple
ExecStart=%h/.local/share/hyprwhspr-ai/.venv/bin/python -m src.daemon
WorkingDirectory=%h/.local/share/hyprwhspr-ai
Restart=on-failure
RestartSec=5
# The DAEMON itself idles at ~50 MB. But it shells out to Surya for
# `--engine surya` and `--engine hybrid` OCR modes — Surya loads PyTorch
# ROCm + recognition models into the daemon's cgroup and needs ~2.5 GB
# RSS. Without headroom it thrashes through swap (1+ GB swapped, OCR
# takes minutes instead of seconds). 8 GB cap keeps a runaway in check
# while letting Surya breathe.
MemoryHigh=4G
MemoryMax=8G
# Inherit the user session env (HYPRLAND_INSTANCE_SIGNATURE, XDG_RUNTIME_DIR).
PassEnvironment=HYPRLAND_INSTANCE_SIGNATURE XDG_RUNTIME_DIR

[Install]
WantedBy=default.target
