# dot_config/systemd/user/hyprwhspr-nllb.service
[Unit]
Description=NLLB-200 translation server (lazy-load, idle-exits after 10min)
Documentation=https://github.com/facebookresearch/fairseq/tree/nllb

[Service]
Type=simple
ExecStart=%h/.local/share/hyprwhspr-translate/.venv/bin/python %h/.local/share/hyprwhspr-translate/nllb-server.py
# Don't auto-restart on idle exit (the client will start us back up on next request).
Restart=no
TimeoutStopSec=10

[Install]
WantedBy=default.target
