# dot_config/systemd/user/waybar.service
[Unit]
Description=Waybar — status bar
Documentation=https://github.com/Alexays/Waybar/wiki
PartOf=graphical-session.target
After=graphical-session.target

[Service]
# Canonical: waybar auto-discovers ~/.config/waybar/config.jsonc.
ExecStart=/usr/bin/waybar
ExecReload=/usr/bin/kill -SIGUSR2 $MAINPID
KillMode=mixed
Restart=on-failure
RestartSec=2
# Cap cascading restarts when something keeps killing waybar (e.g. config
# parse error after an edit) — fail loud after 5 retries in 60s instead of
# spinning forever and burning CPU.
StartLimitBurst=5
StartLimitIntervalSec=60

[Install]
WantedBy=graphical-session.target
