# dot_config/systemd/user/hyprwhspr.service.d/audio-source.conf
[Service]
# Pin capture to the cleaned-mic PipeWire source regardless of which source
# is currently the pulse default. The pulse-compat layer (PortAudio hostapi=0,
# device 'pulse') honors PULSE_SOURCE per-client and routes through PipeWire's
# spa resampler 48kHz -> 16kHz, which is sub-noise quality for ASR mel
# features (8kHz Nyquist). Without this pin, mid-session changes to the
# default source (e.g. screen-share rerouting) would steal hyprwhspr's input.
Environment=PULSE_SOURCE=voice_cleaned
