# dot_config/systemd/user/hyprwhspr.service.d/rewrite-model.conf
[Service]
# Pin the rewrite hook to Gemma 4 E4B regardless of which other models are
# loaded in LM Studio. Set to "auto" to auto-pick the first non-embedding
# model returned by /v1/models.
Environment=HYPRWHSPR_LMS_MODEL=google_gemma-4-e4b-it
