# dot_config/systemd/user/hyprwhspr.service.d/wayland-env.conf
[Service]
# Subprocesses spawned by hyprwhspr (wl-copy, the MIC-OSD daemon, etc.)
# need WAYLAND_DISPLAY to reach the compositor. The user manager doesn't
# always have it imported before this unit starts (race with omarchy's
# import-environment in autostart.conf), so set it explicitly. Hyprland on
# this machine consistently uses wayland-1 — adjust if that ever changes.
Environment=WAYLAND_DISPLAY=wayland-1
