# dot_config/systemd/user/app-walker@autostart.service.d/restart.conf
# Auto-restart the walker --gapplication-service if it dies.
# Mirrors omarchy's default/walker/restart.conf — drop-in for the
# systemd user unit `app-walker@autostart.service` that walker.desktop
# (XDG autostart) creates implicitly.
#
# Activated by enabling the unit:  systemctl --user enable app-walker@autostart.service

[Service]
Restart=always
RestartSec=2
