# dot_config/systemd/user/wayland-wm@hyprland.desktop.service.d/oom.conf
# Protect the Wayland compositor from the kernel OOM killer.
# Default (omarchy/uwsm) ships OOMScoreAdjust=200, which marks the compositor
# as a PREFERRED kill target — backwards. A runaway job killed the whole
# session this way on 2026-06-19 (30GB python3 -> global OOM).
#
# Heavy jobs now run in their own slice (see ~/.local/bin/heavy), so they no
# longer inherit this negative score from the session cgroup.
#
# Takes effect on next compositor start (i.e. next login).
[Service]
OOMScoreAdjust=-800
