# dot_config/systemd/user/hyprwhspr-ai.service.d/no-keepalive.conf
[Service]
# Disable the daemon's internal LM Studio keepalive loop. With this set,
# the daemon never auto-pings LM Studio to refresh the JIT TTL; the rewrite
# model only loads in response to actual transcription/rewrite requests
# and unloads naturally after LM Studio's TTL elapses (default 60min).
# Set to a positive value (e.g. 300) to re-enable periodic keepalive.
Environment=HYPRWHSPR_AI_KEEPALIVE_INTERVAL=0
