# dot_config/systemd/user/hyprwhspr-ai.service.d/wayland-env.conf
[Service]
# Daemon shells out to `grim` (and `slurp`) for screen capture in the
# vision/OCR paths. Both need WAYLAND_DISPLAY to reach the compositor.
# The unit's PassEnvironment includes HYPRLAND_INSTANCE_SIGNATURE +
# XDG_RUNTIME_DIR but not WAYLAND_DISPLAY — set it explicitly to match
# the hyprwhspr.service drop-in (Hyprland on this machine consistently
# uses wayland-1). Without this, `hyprwhspr-ai vision-task` fails with
# "vision-task failed: grim failed" and the menu surfaces an "empty"
# result with no error (stderr was suppressed there pre-fix).
Environment=WAYLAND_DISPLAY=wayland-1
